# systemd/instapaper-delivery-rich-read.service
[Unit]
Description=Weekly Instapaper Delivery rich read
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=kunalgupta
WorkingDirectory=/home/kunalgupta/Desktop/knlgpt/Instapaper_Delivery
Environment=CODEX_HOME=/home/kunalgupta/.codex
Environment="INSTAPAPER_DELIVERY_NAME=Rich Read"
Environment=INSTAPAPER_DELIVERY_FOLDER_ID=5255436
Environment="INSTAPAPER_DELIVERY_FOLDER_NAME=Rich Reads"
Environment="INSTAPAPER_DELIVERY_SUMMARY_LABEL=Rich Read Summary"
EnvironmentFile=/home/kunalgupta/.config/knlgpt-orchestration/instapaper-delivery.env
LoadCredentialEncrypted=instapaper_delivery_credentials
LoadCredentialEncrypted=resend_api_key
ExecStart=/home/kunalgupta/.local/share/instapaper-delivery/run_daily_instapaper_delivery.sh
