# systemd/instapaper-delivery-snap-read.service
[Unit]
Description=Daily Instapaper Delivery snap read
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=kunalgupta
WorkingDirectory=/home/kunalgupta/Desktop/knlgpt/Instapaper_Delivery
Environment=CODEX_HOME=/home/kunalgupta/.codex
EnvironmentFile=/home/kunalgupta/.config/knlgpt-orchestration/instapaper-delivery.env
LoadCredentialEncrypted=instapaper_delivery_credentials
LoadCredentialEncrypted=resend_api_key
ExecStart=/home/kunalgupta/.local/share/instapaper-delivery/run_daily_instapaper_delivery.sh
